# apk/src/digitalocean
#!/sbin/openrc-run

description="Loads Digital Ocean droplet configuration"

required_files="/media/cdrom/digitalocean_meta_data.json"

depend() {
	need localmount
	before network
	before hostname
}

start() {
	ebegin "Loading Digital Ocean configuration"

	hostname="$(jq -jre '.hostname' /media/cdrom/digitalocean_meta_data.json 2>/dev/null)"
	if [ $? -eq 0 ]; then
		echo "$hostname" > /etc/hostname
	else
		ewarn "Could not read hostname"
	fi

	f="/etc/network/interfaces"

	echo "auto lo" > "$f"
	echo "iface lo inet loopback" >> "$f"

	ip_addr="$(jq -jre '.interfaces.public[0].ipv4.ip_address' /media/cdrom/digitalocean_meta_data.json 2>/dev/null)"
	ip_netmask="$(jq -jre '.interfaces.public[0].ipv4.netmask' /media/cdrom/digitalocean_meta_data.json 2>/dev/null)"
	ip_gateway="$(jq -jre '.interfaces.public[0].ipv4.gateway' /media/cdrom/digitalocean_meta_data.json 2>/dev/null)"

	echo >> "$f"
	echo "auto eth0" >> "$f"
	echo "iface eth0 inet static" >> "$f"
	echo "	address $ip_addr" >> "$f"
	echo "	netmask $ip_netmask" >> "$f"
	echo "	gateway $ip_gateway" >> "$f"

	ip_addr="$(jq -jre '.interfaces.public[0].ipv6.ip_address' /media/cdrom/digitalocean_meta_data.json 2>/dev/null)"
	ip_error=$?
	ip_cidr="$(jq -jre '.interfaces.public[0].ipv6.cidr' /media/cdrom/digitalocean_meta_data.json 2>/dev/null)"
	ip_gateway="$(jq -jre '.interfaces.public[0].ipv6.gateway' /media/cdrom/digitalocean_meta_data.json 2>/dev/null)"

	if [ $ip_error -eq 0 -a -n "$ip_addr" ]; then
		modprobe ipv6

		echo >> "$f"
		echo "iface eth0 inet6 static" >> "$f"
		echo "	address $ip_addr" >> "$f"
		echo "	netmask $ip_cidr" >> "$f"
		echo "	gateway $ip_gateway" >> "$f"
		echo "	pre-up echo 0 > /proc/sys/net/ipv6/conf/eth0/accept_ra" >> "$f"
	fi

	ip_addr="$(jq -jre '.interfaces.private[0].ipv4.ip_address' /media/cdrom/digitalocean_meta_data.json 2>/dev/null)"
	ip_error=$?
	ip_netmask="$(jq -jre '.interfaces.private[0].ipv4.netmask' /media/cdrom/digitalocean_meta_data.json 2>/dev/null)"

	if [ $ip_error -eq 0 -a -n "$ip_addr" ]; then
		echo >> "$f"
		echo "auto eth1" >> "$f"
		echo "iface eth1 inet static" >> "$f"
		echo "	address $ip_addr" >> "$f"
		echo "	netmask $ip_netmask" >> "$f"
	fi

	eend 0
}
